# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-apps/sysmocom-backup-data/files/sysmocom-bsc/data.mount
[Unit]
Description=/data

[Mount]
What=/dev/sda4
Where=/data
Type=ext4
