# recipes-support/lc15-sysdev-remap/files/mnt-rom-user.mount
[Unit]
Description=Mount user partition with default settings

[Mount]
What=/dev/mtdblock9
Where=/mnt/rom/user
Type=jffs2

[Install]
WantedBy=multi-user.target
