# recipes-support/lc15-sysdev-remap/files/mnt-rom-factory.automount
[Unit]
Description=Automount ROM partition with factory defaults

[Automount]
Where=/mnt/rom/factory

[Install]
WantedBy=multi-user.target
