# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-support/lc15-sysdev-remap/files/mnt-storage.automount
[Unit]
Description=Automount storage partition

[Automount]
Where=/mnt/storage

[Install]
WantedBy=multi-user.target
