# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-support/lc15-sysdev-remap/files/mnt-rom-factory.mount
[Unit]
Description=Mount ROM partition with factory defaults

[Mount]
What=/dev/mtdblock8
Where=/mnt/rom/factory
Type=jffs2
Options=ro

[Install]
WantedBy=multi-user.target
