# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-support/lc15-sysdev-remap/files/mnt-rom-user.automount
[Unit]
Description=Automount user partition with default settings

[Automount]
Where=/mnt/rom/user

[Install]
WantedBy=multi-user.target
