# recipes-core/busybox/busybox-systemd/busybox-syslog.service.in
[Unit]
Description=System Logging Service
Wants=busybox-klogd.service

[Service]
ExecStart=/bin/sh -c '. /etc/default/busybox-syslog; exec @base_sbindir@/syslogd -n $SYSLOG_ARGS
Sockets=syslog.socket

[Install]
WantedBy=multi-user.target
Also=busybox-klogd.service
Alias=syslog.service
