# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-bsp/ntp/files/ntpdate.cron
*/30 * * * * root __bindir__/ntpdate-cron
