# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-config/sysmocom-systemd/files/emergency.service
[Unit]
Description=Emergency Reboot
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target

[Service]
ExecStart=/bin/systemctl --force reboot
Type=oneshot
