# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-apps/sysmocom-backup-data/files/sysmocom-odu/data.mount
[Unit]
Description=/data

[Mount]
What=ubi0:data
Where=/data
Type=ubifs
