# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-apps/sysmocom-udhcpd-config/files/sysmocom-bsc/udhcpd.service
[Unit]
Description=UDHCPD
Requires=sys-subsystem-net-devices-eth1.device
After=sys-subsystem-net-devices-eth1.device

[Service]
Type=simple
ExecStart=/usr/sbin/udhcpd -fS /etc/udhcpd.conf

[Install]
WantedBy=multi-user.target
