# software/osmo/system-images/git/meta-sysmocom-bsp/recipes-apps/sysmocom-backup-data/files/sysmocom-restore.service
[Unit]
Description=Restore Configuration Backup
Requires=data.mount
After=data.mount

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/sysmocom-restore-data /data/sysmocom-backup.tar

[Install]
WantedBy=multi-user.target
